# .config/autostart/startup.desktop
[Desktop Entry]
Type=Application
Name=Startup
Exec=/home/dean/startup.sh
Comment=Startup script
